# ros2mower.service
[Unit]
Description=ROS2Mower

[Service]
Type=simple
ExecStart=/root/boot_ros2mower.sh start

[Install]
WantedBy=multi-user.target