# ros2mower.service
[Unit]
Description=ROS2Mower

[Service]
Type=simple
ExecStart=/data/ROS2Mower/ros2mower.sh start

[Install]
WantedBy=multi-user.target